# docker/root/etc/s6-overlay/s6-rc.d/svc-devshop/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 80" \
    s6-setuidgid abc \
      /usr/share/devshop/scripts/hosting-queue-runner
