# roles/geerlingguy.supervisor/templates/program.conf.j2
[program:{{ item.name }}]
command={{ item.command }}
{% if item.configuration is defined %}
{{ item.configuration }}
{% endif %}
